# com.newsfeed.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.newsfeed.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/xujintao/Desktop/news-agg/server.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/xujintao/Desktop/news-agg</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/xujintao/Desktop/news-agg/server.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/xujintao/Desktop/news-agg/server.log</string>
</dict>
</plist>
